# root/etc/hotplug.d/iface/01-transparent-proxy
#!/bin/sh

ipset flush tp_spec_def_gw
ipset destory tp_spec_def_gw
ipset create tp_spec_def_gw hash:net hashsize 64
DEFAULT_GATEWAY=$(ip route | grep default | grep -Eoe "via (\d+.\d+.\d+.\d+)" | awk '{print $NF}')
ipset add tp_spec_def_gw $(ip route | grep $DEFAULT_GATEWAY | grep -Eoe "src (\d+.\d+.\d+.\d+)" | awk '{print $NF}')
